# skills/screenwatch/launchd/dashboard.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!-- Dashboard HTTP server. install.sh substitutes __HOME__ and __LABEL__.
     serve.py binds 127.0.0.1:8484 by default; set SCREENWATCH_HOST/PORT in the
     EnvironmentVariables dict below to expose it (no auth — opt in knowingly). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<true/>
	<key>Label</key>
	<string>__LABEL__</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/python3</string>
		<string>__HOME__/screenwatch/bin/serve.py</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>__HOME__/screenwatch/logs/dashboard-server.log</string>
	<key>StandardOutPath</key>
	<string>__HOME__/screenwatch/logs/dashboard-server.log</string>
</dict>
</plist>
